# caddyfile-manager.service
[Unit]
Description=Caddyfile可视化编辑器服务
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/caddyfile-manager
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="CADDYFILE_PATH=/etc/caddy/Caddyfile"
Environment="CADDY_BINARY=caddy"
Environment="PORT=5000"
Environment="HOST=0.0.0.0"
ExecStart=/usr/bin/python3 /opt/caddyfile-manager/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

